# scripts/com.lifesailor.securities-reports.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lifesailor.securities-reports</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/lifesailor/Desktop/kosmos/ai/investment/market_summary/.venv/bin/python</string>
    <string>/Users/lifesailor/Desktop/kosmos/ai/investment/market_summary/scripts/collect_securities_reports.py</string>
  </array>

  <!-- 매주 일요일 19:30 KST -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key> <integer>0</integer>
    <key>Hour</key>    <integer>19</integer>
    <key>Minute</key>  <integer>30</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/lifesailor/Desktop/kosmos/ai/investment/market_summary</string>

  <key>StandardOutPath</key>
  <string>/Users/lifesailor/Desktop/kosmos/ai/investment/market_summary/logs/securities_reports.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/lifesailor/Desktop/kosmos/ai/investment/market_summary/logs/securities_reports.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/lifesailor/Desktop/kosmos/ai/investment/market_summary/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
